# dotfiles_bootstrap_ansible/roles/files_setup/files/usr/lib/systemd/system/vbox-modules-sign-secureboot.service
[Unit]
Description=Sign VirtualBox Modules
After=network.target
Before=vboxdrv.service

[Service]
Type=oneshot
ExecStart=/root/.MOK/vbox-modules-sign-secureboot.sh 

[Install]
WantedBy=default.target

